# bin/service/vsftpd/run
#!/bin/bash

echo "Starting Vsftpd"

exec /usr/sbin/vsftpd
